# deploy/proxypool.service
[Unit]
Description=ProxyPool - Free Proxy Aggregator & Clash Subscription
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 项目目录与 venv（按实际路径修改）
WorkingDirectory=%h/桌面/工作目录/proxypool
ExecStart=%h/桌面/工作目录/proxypool/.venv/bin/python main.py
# 崩溃后 10 秒自动重启
Restart=on-failure
RestartSec=10
# 标准输出/错误走 journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
